# packaging/launchd/com.macdoctor.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macdoctor.scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__NOTIFY_SCRIPT__</string>
    </array>

    <!-- Run every Sunday at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/mac-doctor-scan.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mac-doctor-scan.log</string>

    <!-- Run missed scans when Mac wakes from sleep -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
